# kernel/darqos.target
[Unit]
Description=DarqOS Runlevel
Requires=multi-user.target
